# etc/gniza-daemon.service
[Unit]
Description=GNIZA Background Health Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/gniza/bin/gniza daemon start
WorkingDirectory=/usr/local/gniza
Environment=GNIZA_DIR=/usr/local/gniza
Environment=PYTHONPATH=/usr/local/gniza
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gniza-daemon
NoNewPrivileges=yes
ProtectSystem=no
PrivateTmp=yes
ProtectHome=no
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
