# etc/gniza-web.service
[Unit]
Description=GNIZA Web Dashboard
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m web --host=0.0.0.0 --port=2323
WorkingDirectory=/usr/local/gniza
Environment=GNIZA_DIR=/usr/local/gniza
Environment=PYTHONPATH=/usr/local/gniza
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
