# etc/gniza-web.service
[Unit]
Description=GNIZA Web Dashboard
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m web
WorkingDirectory=/usr/local/gniza
Environment=GNIZA_DIR=/usr/local/gniza
Environment=GNIZA_CONFIG_DIR=/etc/gniza
Environment=LOG_DIR=/var/log/gniza
Environment=PYTHONPATH=/usr/local/gniza
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
